# com.thc.fleetmap.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thc.fleetmap.morning</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/willlawrence/Desktop/Willy/FleetMapAndTimeline/generate.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/fleetmap-morning.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/fleetmap-morning.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/willlawrence/Desktop/Willy/FleetMapAndTimeline</string>
</dict>
</plist>
